# config/includes.chroot/etc/xdg/autostart/seapath-installer.desktop
[Desktop Entry]
Type=Application
Name=seapath-installer
Exec=start-seapath-installer
StartupNotify=false
Terminal=false